# debian/hasura.service
[Unit]
Description=hasura

[Service]
EnvironmentFile=-/etc/hasura/hasura.conf
ExecStart=/usr/bin/hasura serve
